# conf/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:samp03svr]
command=~/server/samp03svr
process_name=samp03svr
autorestart=true
