# conf/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:samp03svr]
user=root
directory=/root/server/
command=./samp03svr
process_name=samp03svr
autostart=true
autorestart=true
redirect_stderr=true
